# deploy/unichat.service
[Unit]
Description=UniChat 消息聚合平台
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=unichat
Group=unichat
WorkingDirectory=/opt/unichat
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

ExecStartPre=/usr/bin/env uv sync --frozen
ExecStart=/usr/bin/env uv run uvicorn src.main:app --host 127.0.0.1 --port 8000

Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

StandardOutput=journal
StandardError=journal

# Security hardening
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
